# supervisor/frontend.conf
[program:frontend]
directory=/app/frontend
command=yarn start
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/frontend.err.log
stdout_logfile=/var/log/supervisor/frontend.out.log
stopsignal=TERM
stopwaitsecs=20
stopasgroup=true
killasgroup=true
environment=HOST="0.0.0.0",PORT="3000"
